# kb-mcp/examples/deployments/intranet-http/kb-mcp.service
# /etc/systemd/system/kb-mcp.service
#
# systemd unit for the intranet HTTP deployment. See README.md / README.ja.md
# for full setup. Enable with:
#     sudo cp kb-mcp.service /etc/systemd/system/
#     sudo systemctl daemon-reload
#     sudo systemctl enable --now kb-mcp.service

[Unit]
Description=kb-mcp Streamable HTTP MCP server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kbmcp
Group=kbmcp

# CWD = /srv/kb-mcp so that Config::discover() picks up
# /srv/kb-mcp/kb-mcp.toml automatically (priority 2 — CWD).
WorkingDirectory=/srv/kb-mcp

# `serve` reads transport / kb_path / model from kb-mcp.toml — keep ExecStart minimal.
ExecStart=/usr/local/bin/kb-mcp serve

# Restart on transient failures (network blips, OOM kills, etc.).
Restart=on-failure
RestartSec=5s

# Resource hardening. Tighten further on production hosts.
# kb-mcp needs no privileged operations, no inbound device access.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
RestrictRealtime=true

# Read-write paths kb-mcp legitimately needs:
#   - /srv/kb-mcp/    : kb_path itself + .kb-mcp.db (the parent of kb_path)
#   - /var/cache/fastembed : ONNX model cache (matches FASTEMBED_CACHE_DIR in kb-mcp.toml)
ReadWritePaths=/srv/kb-mcp /var/cache/fastembed

# RUST_LOG controls verbosity. info gives the "loaded config source=..." startup
# line plus warnings; bump to debug for troubleshooting.
Environment=RUST_LOG=info

# stdout/stderr go to the journal — `journalctl -u kb-mcp.service -f` to tail.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
